# deploy/claudeclaw.service
# ClaudeClaw systemd user service
# Install: deploy.sh copies this to ~/.config/systemd/user/ and templates paths
# Do not use %h directly — deploy.sh replaces __HOME__ and __PROJECT__ at install time

[Unit]
Description=ClaudeClaw Telegram Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=__PROJECT__
EnvironmentFile=__PROJECT__/.env
Environment=PATH=__HOME__/.nvm/versions/node/v22.22.2/bin:__HOME__/.npm-global/bin:__HOME__/.local/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=__BUN__ run src/index.ts
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Shutdown — bound the SIGTERM wait so a stuck SDK Stop hook (codex plugin's
# stop-review-gate-hook has a 900s timeout) cannot pin systemd for 90s.
# KillMode=mixed: SIGTERM to MainPID; if cleanup hangs past TimeoutStopSec
# systemd sends SIGKILL to all remaining cgroup processes (works around the
# WSL2 cgroup-kill bug we hit on 2026-05-01).
TimeoutStopSec=15
KillMode=mixed
KillSignal=SIGTERM

# Security — allowlist is the boundary, not filesystem sandbox
NoNewPrivileges=true
PrivateNetwork=false
MemoryDenyWriteExecute=false

[Install]
WantedBy=default.target
