# deploy/claudeclaw.service
# ClaudeClaw systemd user service
# Install: deploy.sh copies this to ~/.config/systemd/user/ and templates paths
# Do not use %h directly — deploy.sh replaces __HOME__ and __PROJECT__ at install time

[Unit]
Description=ClaudeClaw Telegram Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=__PROJECT__
EnvironmentFile=__PROJECT__/.env
Environment=PATH=__HOME__/.npm-global/bin:__HOME__/.local/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=__BUN__ run src/index.ts
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security — allowlist is the boundary, not filesystem sandbox
NoNewPrivileges=true
PrivateNetwork=false
MemoryDenyWriteExecute=false

[Install]
WantedBy=default.target
